# playbooks/roles/private-node/000_autossh/templates/autossh-tunnel.service.j2
[Unit]
Description=AutoSSH tunnel to VPS
After=network.target

[Service]
User={{ linux_user }}
ExecStart=/usr/bin/autossh -M 0 -N -o "ServerAliveInterval 60" -o "ServerAliveCountMax 30" -o "StrictHostKeyChecking=accept-new" -i /home/{{ linux_user }}/.ssh/{{ autossh_ssh_key_file }} -R 0.0.0.0:{{ master_port }}:localhost:22 {{ master_linux_user }}@{{ master_ip }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
